# precise-mongodb/config/supervisor/conf.d/mongodb.conf
[program:mongodb]
command=/bin/sh /etc/init.d/mongod start